# schedule.cron
# mishearing-bot stuff
# 40 * * * * cd /var/www/not-an-alien-bot && make run-mishear-fact
10,20,25,30,40 * * * * cd /var/www/not-an-alien-bot && make run-mishear-quote
