# schedule.cron
# mishearing-bot stuff
# 40 * * * * cd /var/www/not-an-alien-bot && make run-mishear-fact
56,26 * * * * cd /var/www/not-an-alien-bot && make run-mishear-quote
