# schedule.cron
# mishearing-bot stuff
15,50,56,26 * * * * cd /usr/src/app && make run-mishear-quote
# * * * * * echo "Hello world" >> /var/log/cron.log 2>&1

